# video-transcription.service
[Unit]
Description=YouTube Video Transcription Service
After=network.target

[Service]
Type=simple
User=rgann
WorkingDirectory=/home/rgann/video-processing-service
ExecStart=/home/rgann/video-processing-service/venv/bin/python cloud-vm-processor.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=video-transcription
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/rgann/video-processing-service/.env

[Install]
WantedBy=multi-user.target