# video-transcription.service
[Unit]
Description=YouTube Video Transcription Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/youtube-video-transcriptor
ExecStart=/home/ubuntu/youtube-video-transcriptor/venv/bin/python cloud-vm-processor.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=video-transcription
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/ubuntu/youtube-video-transcriptor/.env

[Install]
WantedBy=multi-user.target